# packages/sing-box-config/sing_box_config-0.2.1.tar.gz/sing_box_config-0.2.1/roles/sing_box_tproxy/templates/etc/systemd/system/sing-box-reload.service.j2
[Unit]
Description=Reload sing-box after config changes

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reload sing-box.service
